# dots/autostart/Tint2.desktop
[Desktop Entry]
Exec=tint2
Name=Tint2
Type=Application
Version=1.0
